# tests/service/service1/run
#!/bin/sh

echo '{"test": "foo"}'
echo "test" >&2
echo foo

sleep 10
