# tests/service/service4/run
#!/bin/sh

while true; do
  LC_ALL=C date
  sleep 1
done