# tests/service/service4/run
#!/bin/sh

echo "stderr" >&2

while true; do
  LC_ALL=C date
  sleep 1
done