# tests/service/service2/run
#!/bin/sh

echo '{"so": "me"}'

(bash -c "sleep 3; echo x y z" &)

sleep infinity
